# scripts/com.richardcarroll.kit-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.richardcarroll.kit-sync</string>

    <!--
        The Notion token is NOT stored here. kit_sync_launch.sh loads it from
        an external secrets file (~/.config/kit_sync/notion_token) and exports
        NOTION_TOKEN before exec'ing the runner. See scripts/kit_sync_launch.sh.
    -->
    <key>ProgramArguments</key>
    <array>
        <string>/Users/richardcarroll/Library/Scripts/kit_sync_launch.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/richardcarroll/Library/Scripts</string>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>5</integer>
        <key>Minute</key>
        <integer>45</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/richardcarroll/Library/Logs/kit-sync.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/richardcarroll/Library/Logs/kit-sync.log</string>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
